# scheduler/org.cups.cups-lpd.socket
[Unit]
Description=CUPS LPD Server Socket

[Socket]
ListenStream=515
Accept=yes

[Install]
WantedBy=sockets.target
