# scheduler/org.cups.cupsd.service.in
[Unit]
Description=CUPS Scheduler
Documentation=man:cupsd(8)
After=sssd.service

[Service]
ExecStart=@sbindir@/cupsd -l
Type=simple
Restart=on-failure

[Install]
Also=org.cups.cupsd.socket org.cups.cupsd.path
WantedBy=printer.target
