# scheduler/org.cups.cupsd.service.in
[Unit]
Description=CUPS Scheduler

[Service]
ExecStart=@sbindir@/cupsd -l
Type=simple

[Install]
Also=org.cups.cupsd.socket org.cups.cupsd.path
WantedBy=printer.target
