# scheduler/org.cups.cupsd.path.in
[Unit]
Description=CUPS Scheduler

[Path]
PathExists=@CUPS_CACHEDIR@/org.cups.cupsd

[Install]
WantedBy=multi-user.target
