# scheduler/org.cups.cupsd.service.in
[Unit]
Description=CUPS Scheduler
Documentation=man:cupsd(8)
After=network.target nss-user-lookup.target nslcd.service
Requires=org.cups.cupsd.socket

[Service]
ExecStart=@sbindir@/cupsd -l
Type=notify
Restart=on-failure

[Install]
Also=org.cups.cupsd.socket org.cups.cupsd.path
WantedBy=@SYSTEMD_WANTED_BY@
