# scheduler/org.cups.cupsd.socket.in
[Unit]
Description=CUPS Scheduler
PartOf=org.cups.cupsd.service

[Socket]
ListenStream=@CUPS_DEFAULT_DOMAINSOCKET@

[Install]
WantedBy=sockets.target
